# configs/lightsoff.service
[Unit]
Description=Turns RGB lights off on memory kit
[Service]
ExecStart=/usr/bin/LIGHTSOFF.sh
[Install]
WantedBy=multi-user.target